# systemd/repeat.service
[Unit]
Description=renew let's encrypt

[Service]
Type=oneshot
ExecStart=/home/ubuntu/systemd/repeat.sh

[Install]
WantedBy=default.target
